# workspace/crontab/root
* * * * * laradock /usr/bin/php /var/www/giantechMp/artisan schedule:run >> /dev/null 2>&1
